# roles/branding/files/etc/update-motd.d/20-logo
#!/bin/sh
echo
cat /etc/moodlebox-logo.txt
