# roles/branding/templates/update-motd.d/20-logo
#!/bin/sh
echo
cat /etc/{{ moodlebox_logo_file }}
